# launchd/ai.human-ai.pulse-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.human-ai.pulse-refresh</string>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>21</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- 不再写死 4321:Tauri 壳用动态端口。gw-cron.sh 现场读 ~/.human-ai/.gateway-port。 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/claudecodedezhuanshumac/.human-ai/bin/gw-cron.sh</string>
    <string>/api/pulse/refresh-mirror</string>
    <string>{}</string>
    <string>60</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/claudecodedezhuanshumac/.human-ai/logs/pulse-refresh.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/claudecodedezhuanshumac/.human-ai/logs/pulse-refresh.err.log</string>
</dict>
</plist>
